# scripts/com.xhd.blog.autopost-2am.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xhd.blog.autopost-2am</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/xuhongduo/Documents/myblog/scripts/auto_post.py</string>
    <string>--count</string>
    <string>6</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/xuhongduo/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/xuhongduo</string>
  </dict>

  <!-- 凌晨 2:00 触发，每次生成 6 篇 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>   <integer>2</integer>
    <key>Minute</key> <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>       <false/>
  <key>LaunchOnlyOnce</key>  <false/>

  <key>WorkingDirectory</key>
  <string>/Users/xuhongduo/Documents/myblog</string>

  <key>StandardOutPath</key>
  <string>/Users/xuhongduo/Documents/myblog/scripts/logs/launchd-2am.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/xuhongduo/Documents/myblog/scripts/logs/launchd-2am-err.log</string>
</dict>
</plist>
